# deploy/udev/70-wlan0-powersave-off.rules
ACTION=="add|change", SUBSYSTEM=="net", KERNEL=="wlan0", RUN+="/bin/sh -c 'iw dev wlan0 set power_save off'"
