# deploy/systemd/ess-thermal-check.service

[Unit]
Description=ESS Thermal (MLX90640) Presence Check
After=systemd-modules-load.service
Wants=systemd-modules-load.service
Before=ess-vision.service
PartOf=ess.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/ess-guardian/current/scripts/ess-thermal-check.sh
RemainAfterExit=yes
Restart=on-failure
RestartSec=2

[Install]
WantedBy=ess.target
